# templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
Type=simple
ExecStart={{ exec_node_exporter }}

[Install]
WantedBy=multi-user.target
